# projects/S905/filesystem/usr/lib/systemd/system/amlogic-wifi.service
[Unit]
Description=Amlogic WiFi module loader

[Service]
Type=simple
ExecStart=/sbin/modprobe wifi_dummy

[Install]
WantedBy=basic.target
