# projects/S905/packages/linux-firmware/brcmfmac_sdio-firmware-aml/system.d/brcmfmac_sdio-firmware-aml.service
[Unit]
Description=Broadcom sdio firmware update for %I
ConditionPathExists=/dev/ttyS1
After=network.target

[Service]
Type=simple
RemainAfterExit=yes
ExecStart=/usr/bin/hciattach -n -s 115200 /dev/ttyS1 bcm43xx 2000000

[Install]
WantedBy=network.target
