# projects/S905/packages/linux-firmware/brcmfmac_sdio-firmware-aml/system.d/brcmfmac_sdio-firmware@.service
[Unit]
Description=Broadcom sdio firmware update for %I
ConditionPathExists=/dev/ttyS1

[Service]
Type=simple
RemainAfterExit=yes
ExecStart=/usr/bin/brcm_patchram_plus --patchram /lib/firmware/brcm/%I.hcd /dev/ttyS1 --baudrate 2000000 --enable_hci --no2bytes --tosleep=50000
