# projects/S905_1GB/filesystem/usr/lib/systemd/system/fstrim.service
[Unit]
Description=fstrim storage partition

[Service]
Type=simple
ExecStart=-/usr/sbin/fstrim -a -v

[Install]
WantedBy=basic.target
