# projects/S905/packages/linux-firmware/qca9377-firmware-aml/system.d/qca9377-firmware-aml.service
[Unit]
Description=QCA Bluetooth firmware service
ConditionPathExists=/dev/ttyS1

[Service]
Type=simple
RemainAfterExit=yes
ExecStart=/usr/bin/hciattach -n -s 2000000 /dev/ttyS1 qca
